# src/etc/systemd/system/notify_reboot.timer
[Unit]
Description=send email on reboot

[Timer]
OnBootSec=60
Unit=notify_reboot.service

[Install]
WantedBy=multi-user.target
